# build/virtualbox/etc/systemd/system/serverjockey-welcome.service
[Unit]
Description=ServerJockey Welcome Banner
Requires=sysinit.target

[Service]
Type=oneshot
ExecStart=/root/welcome.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
